# supervisor/gunicorns.conf
[program:gunicorn]
command=gunicorn --pythonpath /home/app/web config.wsgi:application --bind 0.0.0.0:8000
directory=/home/app/web
autorestart=true
redirect_stderr=true
